# Stellar/horizon.service
[Unit]
Description=horizon
 
[Service]
Type=simple
User=ubuntu
Restart=always
ExecStart=/usr/bin/horizon
 
[Install]
WantedBy=default.target